# pi-airplay.service
[Unit]
Description=Pi-AirPlay Complete Service
After=network.target sound.target
# No longer requires shairport-sync.service since we'll handle it ourselves

[Service]
ExecStart=/bin/bash /home/ivpi/pi-airplay/start_pi_airplay.sh
Restart=always
RestartSec=5
User=ivpi
Group=ivpi
WorkingDirectory=/home/ivpi/pi-airplay
# Standard service settings
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target