# config/airplay-web.service
[Unit]
Description=AirPlay Web Interface
After=network.target shairport-sync.service

[Service]
ExecStart=/usr/bin/python3 /home/pi/airplay-receiver/app.py
WorkingDirectory=/home/pi/airplay-receiver
User=pi
Environment=FLASK_ENV=production
Environment=FLASK_APP=app.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/airplay-web.log
StandardError=append:/var/log/airplay-web.log

[Install]
WantedBy=multi-user.target