# config/music-display.service
[Unit]
Description=Music Display with AirPlay and Recognition
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/music-display
# Using a bash script wrapper that handles virtual environment activation
ExecStart=/bin/bash /home/pi/music-display/start_music_display.sh
Restart=always
StandardOutput=journal
StandardError=journal
SyslogIdentifier=music-display
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target